# setup/systemd/radar-update-check.service
[Unit]
Description=Radar Project - Update Checker (runs at startup)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
ExecStart=/opt/radar-update/check-updates.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
