# setup/systemd/radar-server.service
[Unit]
Description=Radar Full-Stack Server
Documentation=https://github.com/yourusername/radar-project
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/radar-project/RaspberryPiRadarFullStackApplicationAndStepperController/server
EnvironmentFile=/home/pi/radar-project/RaspberryPiRadarFullStackApplicationAndStepperController/server/.env
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="NODE_ENV=production"

# Start the server
ExecStart=/usr/bin/node index.js

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radar-server

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
